# exporter/src/main/etc/supervisor/conf.d/exporter.conf
[program:exporter]
directory=/full/path/exporter
command=java -jar exporter-1.0-SNAPSHOT.jar
autostart=true
autorestart=true
redirect_stderr=true
user=user
